# config/mise/tasks.d/home-server/systemd/user/yt-dlp-update.service
[Unit]
Description=Upgrade yt-dlp with mise

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mise upgrade yt-dlp
