# config/mise/tasks.d/home-server/systemd/user/yt-dlp-update.timer
[Unit]
Description=Run mise upgrade for yt-dlp daily at 9 AM

[Timer]
OnCalendar=*-*-* 09:00
Persistent=true

[Install]
WantedBy=timers.target
